# systemd/papaya-ssh.service
[Unit]
Description=PapayaMeter reverse SSH tunnel service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nvidia
WorkingDirectory=/home/nvidia/Desktop/PapayaMeter
ExecStart=/usr/bin/python3 /home/nvidia/Desktop/PapayaMeter/services/ssh_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=papaya-ssh

[Install]
WantedBy=multi-user.target
